# target/linux/mediatek/filogic/base-files/etc/init.d/turn_on_usb_power
#!/bin/sh /etc/rc.common

START=99

start() {
	case "$(board_name)" in
	glinet,gl-mt3600be)
		echo 1 > /sys/class/gpio/usb_power/value
		;;
	esac
}
